# serialstarter-disable/services/mk2-dbus.ttyS3/run
#!/bin/sh
echo "*** starting mk2-dbus ***"
exec 2>&1
exec softlimit -d 100000000 -s 1000000 -a 100000000 /opt/victronenergy/mk2-dbus/start-mk2-dbus.sh ttyS3 builtin-mkx