# serialstarter-disable/services/vedirect-interface.ttyS1/run
#!/bin/sh
echo "*** starting vedirect-interface ***"
exec 2>&1
exec /opt/victronenergy/vedirect-interface/start-vedirect.sh ttyS1 PRODUCT=builtin-vedirect